# ansible/roles/sshtunnel.upstart.conf
description "SSH Tunnel back to host"
author      "Joshua Eckroth"

setuid vagrant

start on started network
stop on stopping network
stop on starting shutdown

console output

exec ssh -N -L3142:127.0.0.1:3142 cinf401@192.168.121.1
respawn

